# etc/udev/rules.d/120-make-backup.rules
ACTION=="add", KERNEL=="sd*1", RUN+="/opt/make_backup/make_backup_d.sh UDEV"
