# etc/udev/rules.d/120-make-backup.rules
ACTION=="add", KERNEL=="sd*1", RUN+="/etc/make_backup/Make_Backup_d.sh UDEV"
